# deploy/systemd/iot-platform-backend.service
[Unit]
Description=IoT Platform Backend Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=iot-platform
Group=iot-platform
WorkingDirectory=/opt/iot-platform
ExecStart=/usr/bin/java -Xms512m -Xmx2g -XX:+UseG1GC -XX:+UseStringDeduplication -Dspring.profiles.active=prod -jar iot-platform-backend-0.0.1-SNAPSHOT.jar
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=iot-platform-backend

# Environment variables
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=JAVA_OPTS=-Xms512m -Xmx2g -XX:+UseG1GC -XX:+UseStringDeduplication
Environment=DB_USERNAME=postgres
Environment=DB_PASSWORD=123
Environment=JWT_SECRET=shiftAIOT_secure_jwt_secret_key_2024_should_be_at_least_256_bits_long_for_security

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/iot-platform /var/log/iot-platform

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
